# solidpostal.service
[Unit]
Description=solidpostal jenkins configuration

[Service]
ExecStart=/usr/local/sbin/solidpostal

[Install]
WantedBy=multi-user.target